# deploy/systemd/klicklocal-frontend.service
# Runs the PRODUCTION Next.js server (npm start) as a managed service.
# App dir: /var/www/klicklocal/frontend   Port: 3000   Branch: main
#
# Install:
#   sudo cp deploy/systemd/klicklocal-frontend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now klicklocal-frontend
#   sudo systemctl status klicklocal-frontend
#
# Logs:  sudo journalctl -u klicklocal-frontend -f
#
# NOTE: confirm the node path with `which node` and fix ExecStart if needed.

[Unit]
Description=Klicklocal Next.js frontend (production)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/klicklocal/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
