# deploy/supervisor/klicklocal-worker.conf
# Laravel queue worker (scheduled posts / publishing) managed by Supervisor.
#
# Install:
#   sudo cp deploy/supervisor/klicklocal-worker.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl status
#
# NOTE: confirm the php path with `which php` and fix `command` if needed.

[program:klicklocal-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/klicklocal-api/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/klicklocal-api/backend/storage/logs/worker.log
stopwaitsecs=3600
