# deploy/supervisor/klicklocal-worker-staging.conf
# Laravel queue worker — STAGING (scheduled posts / publishing).
# App dir: /var/www/klicklocal-staging/backend   Branch: develop
#
# Install:
#   sudo cp deploy/supervisor/klicklocal-worker-staging.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl status

[program:klicklocal-worker-staging]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/klicklocal-staging/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/klicklocal-staging/backend/storage/logs/worker.log
stopwaitsecs=3600
